# deploy/systemd/trivia-marketing.service
[Unit]
Description=Trivia Engine Marketing Site
Documentation=https://github.com/yourusername/trivia-engine
After=network.target
Wants=trivia-api.service

[Service]
Type=simple
User=trivia
Group=trivia
WorkingDirectory=/home/trivia/trivia-engine/marketing

# Environment variables
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
EnvironmentFile=-/home/trivia/trivia-engine/marketing/.env.local

# Process management
ExecStart=/usr/bin/npm start
ExecReload=/bin/kill -USR2 $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/trivia/trivia-engine/marketing
ReadWritePaths=/home/trivia/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trivia-marketing

[Install]
WantedBy=multi-user.target